# autologin/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin lyh970817 --noclear %I 38400 linux
Type=simple
